# tools/bin/services/api.service
[Unit]
Description=Gpio API
After=network.target

[Service]
ExecStart=/home/pi/Sites/gpio-server/api/bin/www
WorkingDirectory=/home/pi/Sites/gpio-server/api

[Install]
WantedBy=multi-user.target
